# meta-opendreambox/recipes-dreambox/satip-client/satip-client/satip-client.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/bin/satipclient
Restart=on-failure
NotifyAccess=main

[Install]
WantedBy=multi-user.target
