# meta-opendreambox/recipes-dreambox/enigma2/enigma2-streamproxy/streamproxy@.service
[Unit]

[Service]
ExecStart=-/usr/bin/streamproxy
StandardInput=socket
